# ops/systemd/user/codex-wechat-bridge.service
[Unit]
Description=Codex WeChat Bridge
After=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/dev/codex-wechat-bridge
Environment=PATH=%h/.local/share/mise/installs/node/24.13.1/bin:%h/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=CODEX_WECHAT_BRIDGE_CODEX_BIN=%h/.local/share/mise/installs/node/24.13.1/bin/codex
ExecStart=%h/dev/codex-wechat-bridge/.venv/bin/codex-wechat-bridge run
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
